# server/deploy/maruaudio-websocket-auth.service
[Unit]
Description=MaruAudio WebSocket Server
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/www/wwwroot/auth.wzagent.cn/websocket
ExecStart=/www/server/php/82/bin/php server.php
Restart=always
RestartSec=5
StandardOutput=append:/var/log/maruaudio-websocket.log
StandardError=append:/var/log/maruaudio-websocket.log
Environment=WS_PORT=8080
Environment=INTERNAL_PORT=8081
Environment=MARUAUDIO_DB_HOST=localhost
Environment=MARUAUDIO_DB_NAME=maruaudio
Environment=MARUAUDIO_DB_USER=maruaudio
Environment=MARUAUDIO_DB_PASSWORD=6FS64ybEGeyMcpZs
Environment=MARUAUDIO_JWT_SECRET=f701a3e4daf0fb6c80c16d595ee4ba6490a5d47895bc7246c0ef484e77f03326

[Install]
WantedBy=multi-user.target
